# docker/crontab
20 * * * * root cd /opt/check-your-reps && rake congress:check_for_updates >>/var/log/cron.log 2>&1
